# bin/council_bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sovereignshadow.council-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/LegacySafe/Shadow-3-Legacy-Loop-Platform/bin/council_bridge_launcher.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/LegacySafe/Shadow-3-Legacy-Loop-Platform/neural_hub/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Volumes/LegacySafe/Shadow-3-Legacy-Loop-Platform/logs/council_bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/LegacySafe/Shadow-3-Legacy-Loop-Platform/logs/council_bridge_error.log</string>
</dict>
</plist>
